# phase3/v5.5/cron/com.hermes.weekly-memory-synthesis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.weekly-memory-synthesis</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HERMEM_HOME__/v5.5/cron/run_weekly_synthesis.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/hermem-weekly-synthesis.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/hermem-weekly-synthesis.err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
